# snippets/systemd/trigger-smb.mounts.service
[Unit]
Description=Trigger all SMB automounts (boot stabilization)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/trigger-smb-automounts.sh

[Install]
WantedBy=multi-user.target
