# ansible/roles/calibre-importer/templates/calibre-import.service.j2
[Unit]
Description=Calibre auto-import (drop folder -> calibredb)

[Service]
Type=oneshot
ExecStart={{ calibre_importer_script_dest }}
