# ansible/roles/node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
User=node_exporter
{% set textfile_flag = '' %}{% if node_exporter_textfile_dir %}{% set textfile_flag = ' --collector.textfile.directory=' + node_exporter_textfile_dir %}{% endif %}
ExecStart=/usr/local/bin/node_exporter --web.listen-address={{ ansible_host }}:{{ node_exporter_port }}{{ textfile_flag }}
Restart=on-failure
RestartSec=15s

[Install]
WantedBy=multi-user.target
